# ansible-ollama.service
[Unit]
Description=Ansible Ollama Web Service
After=network.target

[Service]
Type=simple
Environment="OPENROUTER_API_KEY="
Environment="OPENROUTER_MODEL=anthropic/claude-3.5-sonnet"
Environment="OPENROUTER_URL=https://openrouter.ai/api/v1/chat/completions"
Environment="USE_CLOUD_FALLBACK=false"
ExecStart=/export/ollama/bin/python3 /usr/local/bin/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
